# services/dropbox_monitor.service
[Unit]
Description=Monitor any changes in 'drop' directory and sync to dropbox
After=i3-gaps.target


[Service]
Type=simple
ExecStart=/bin/bash -c "fswatch --event MovedTo --event MovedFrom --event Renamed --event Removed --event Updated --event Created -o -r %h/drop | xargs -i sh -c 'rclone sync %h/drop/ drop:'"
RemainAfterExit=no

[Install]
WantedBy=i3-gaps.target

